# shairport/rootfs/etc/services.d/shairport-sync/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Shairport Sync with a runtime-generated configuration file.
# ==============================================================================

readonly CONFIG_TEMPLATE="/etc/shairport-sync.conf"
readonly CONFIG_RUNTIME="/etc/shairport-sync.runtime.conf"

bashio::log.info "Starting Shairport Sync add-on..."

# Read add-on options with bashio::config
readonly name="$(bashio::config 'name')"
readonly output_device="$(bashio::config 'output_device')"
readonly output_rate="$(bashio::config 'output_rate')"
readonly output_format="$(bashio::config 'output_format')"
readonly buffer_length_seconds="$(bashio::config 'buffer_length_seconds')"
readonly period_size="$(bashio::config 'period_size')"
readonly buffer_size="$(bashio::config 'buffer_size')"
readonly disable_standby_mode="$(bashio::config 'disable_standby_mode')"
readonly use_precision_timing="$(bashio::config 'use_precision_timing')"
readonly interpolation="$(bashio::config 'interpolation')"
readonly drift_tolerance_seconds="$(bashio::config 'drift_tolerance_seconds')"
readonly ignore_volume_control="$(bashio::config 'ignore_volume_control')"
readonly log_verbosity="$(bashio::config 'log_verbosity')"

if bashio::config.has_value 'mixer_control_name'; then
  readonly mixer_control_name="$(bashio::config 'mixer_control_name')"
else
  readonly mixer_control_name=""
fi

bashio::log.info "AirPlay name: ${name}"
bashio::log.info "ALSA device: ${output_device} (${output_rate} Hz, ${output_format})"
bashio::log.info "Buffer: ${buffer_length_seconds}s, period: ${period_size}, buffer: ${buffer_size}"
bashio::log.info "Standby: ${disable_standby_mode}, precision timing: ${use_precision_timing}, interpolation: ${interpolation}"

if [[ -n "${mixer_control_name}" ]]; then
  bashio::log.info "Hardware mixer: ${mixer_control_name}"
fi

bashio::log.info "Generating Shairport Sync configuration from template..."
tempio -conf /data/options.json -template "${CONFIG_TEMPLATE}" -out "${CONFIG_RUNTIME}"

bashio::log.info "Waiting for NQPTP..."
for _ in $(seq 1 50); do
  if pgrep -x nqptp > /dev/null 2>&1; then
    break
  fi
  sleep 0.1
done

if ! pgrep -x nqptp > /dev/null 2>&1; then
  bashio::log.error "NQPTP is not running; AirPlay 2 requires it on ports 319/320"
  exit 1
fi

bashio::log.info "Launching Shairport Sync..."
exec shairport-sync -c "${CONFIG_RUNTIME}"
